# build/crontab
*/5 * * * * php protected/t4.php send/mail